# DEBIAN/cron.d
#
# Launches influx-gw periodically
#
* * * * * sudo python3.6 /usr/local/share/influx-gw/influx_gw.py /usr/local/share/influx-gw/akom.json